# claude-personality.timer
[Unit]
Description=Run Claude personality generator daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
